# services/vnstat-dashboard.service
[Unit]
Description=VnStat Dashboard Node.js Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/vnstat-web-dashboard
ExecStart=/usr/bin/node server-with-vnstat.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
